# vesteda/crontab
# run python script 
1 * * * * root python /app/vesteda/main.py
