# vesteda/crontab
# run python script 
1 * * * * /usr/local/bin/python /app/vesteda/main.py
